# hestia/tools/one_shots/com.local.hass.mount.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.hass.mount</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/evertappels/hass/hestia/tools/one_shots/hass_mount_once.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>NetworkState</key>
    <true/>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/evertappels/Library/Logs/hass-mount.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/evertappels/Library/Logs/hass-mount.log</string>
</dict>
</plist>
